# roles/youtube/files/youtube-cleanup.service
[Unit]
Description=Cleanup Watched YouTube Videos

[Service]
ExecStart=/usr/local/bin/youtube-cleanup
