# roles/youtube/files/youtube-download.service
[Unit]
Description=Download YouTube Videos

[Service]
ExecStart=/usr/local/bin/youtube-download
